# containers/quadlet/webshot.service
[Unit]
Description=Webshot service (host)
Requires=webshot-postgres.service
Requires=webshot-seaweed.service
After=webshot-postgres.service
After=webshot-seaweed.service
PartOf=webshot-debug-stack.target

[Service]
Type=simple
Environment=CPU_LIMIT=4c
ExecStart=/bin/sh -c 'set -eu; cd "${WEBSHOT_ROOT}"; export LD_LIBRARY_PATH="${WEBSHOT_RUNTIME_LD_LIBRARY_PATH}"; exec "${WEBSHOT_BUILD_DIR}/webshot" --config "${WEBSHOT_ROOT}/config/userver_config.yaml" --config_vars "${WEBSHOT_ROOT}/config/config_vars.debug.yaml"'
Restart=on-failure
