# containers/quadlet/webshot-prod-debug-stack.target
[Unit]
Description=Webshot prod-like debug stack (Quadlet + host service)
Requires=webshot-prod-stack.target
Requires=webshot-prod.service
Requires=webshot-scalar.service
Requires=webshot-reverse-proxy.service
After=webshot-prod-stack.target

