# containers/quadlet/webshot-scalar.container
[Unit]
Description=Webshot Scalar API reference viewer
After=webshot-crawler-network.service
Requires=webshot-crawler-network.service
PartOf=webshot-debug-stack.target

[Service]
Restart=on-failure

[Container]
ContainerName=webshot_scalar
Image=docker.io/scalarapi/api-reference:latest
Network=crawler_net
# JSON config mirrors testing-compose.yaml.
Environment="API_REFERENCE_CONFIG={\"sources\":[{\"url\":\"/openapi/webshot.yaml\"}],\"layout\":{\"hideDownloadButton\":false}}"
Volume=${WEBSHOT_ROOT}/schemas:/openapi:ro
