# containers/quadlet/webshot-test-target.container
[Unit]
Description=Webshot nginx test target on TEST-NET-1
After=webshot-crawler-network.service
Requires=webshot-crawler-network.service
PartOf=webshot-stack.target
PartOf=webshot-prod-stack.target

[Service]
Restart=on-failure

[Container]
ContainerName=webshot_test_target
Image=docker.io/nginx:alpine
Network=crawler_net
NetworkAlias=test.example.com
IP=192.0.2.10
Volume=${WEBSHOT_ROOT}/containers/nginx-test.conf:/etc/nginx/conf.d/default.conf:ro
