# devenv/systemd/webshot.service.in
[Unit]
Description=webshot stack
AssertPathExists=/etc/webshot/config_vars.yaml

[Service]
Type=oneshot
RemainAfterExit=yes
Delegate=cpu memory
DynamicUser=yes
StateDirectory=webshot
WorkingDirectory=@out@
Environment=PATH=@runtimePath@
ExecStart=@repoPython@/bin/python3 -m s6.runtime up --mode prod --layout-root @out@ --state-dir %S/webshot --binary-path @out@/webshotd/webshotd_wrapper --config-vars-source /etc/webshot/config_vars.yaml
ExecStop=@repoPython@/bin/python3 -m s6.runtime down --mode prod --layout-root @out@ --state-dir %S/webshot
TimeoutStartSec=120
TimeoutStopSec=120
KillMode=control-group

[Install]
WantedBy=multi-user.target
