# devenv/systemd/webshot.service.in
[Unit]
Description=@description@
AssertPathExists=/etc/webshot/config_vars.yaml
Conflicts=@conflicts@
@localS3Assert@

[Service]
Type=oneshot
RemainAfterExit=yes
Delegate=cpu memory
DynamicUser=yes
StateDirectory=webshot
RuntimeDirectory=webshot
WorkingDirectory=@out@
ExecStart=@out@/bin/webshot-systemd-runtime up --mode prod --layout-root @out@ --persistent-state-dir %S/webshot --ephemeral-state-dir %t/webshot --binary-path @out@/webshotd/webshotd_wrapper --config-vars-source /etc/webshot/config_vars.yaml@execStartExtra@
ExecStop=@out@/bin/webshot-systemd-runtime down --mode prod --layout-root @out@ --persistent-state-dir %S/webshot --ephemeral-state-dir %t/webshot
TimeoutStartSec=120
TimeoutStopSec=120
KillMode=control-group

[Install]
WantedBy=multi-user.target
